# scheduler/container/template/config/crontab
# DO NOT EDIT THIS FILE MANUALLY!
# THIS FILE IS GENERATED BASED ON THE CONTENTS OF THE jobs/ FOLDER
# ANY MANUAL CHANGES WILL BE OVERWRITTEN!

@reboot quackjob /usr/share/bin/quackjob update-cron-schedule